# rootfs/etc/s6/znc/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

STARTCMD="/usr/bin/znc -f -d /var/lib/znc"

pushd /var/lib/znc > /dev/null
exec /usr/bin/gosu znc ${STARTCMD}
popd > /dev/null
